# mrpir.service
[Unit]
Description=Presense sensor
After=multi-user.target
Wants=network-online.target systemd-networkd-wait-online.service

StartLimitIntervalSec=0
StartLimitBurst=0

[Service]
Type=simple
Restart=on-failure
RestartSec=5s
User=<user>

ExecStart=/usr/bin/python3 /home/<user>/mrpir/mrpir.py

[Install]
WantedBy=multi-user.target
